# config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
port=0.0.0.0:8081

[supervisord]
serverurl=http://127.0.0.1:8081
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
user=root
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:init]
priority=1
command=/opt/setup/scripts/start.sh
autostart=true
startsecs=5
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:redis]
priority=5
command=/usr/bin/redis-server /etc/redis/redis-openvas.conf --daemonize no --syslog-enabled no
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
user=redis
autorestart=true
autostart=false
startretries=10
startsecs=10
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:sshd]
priority=5
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
startsecs=10
autorestart=true
autostart=false
stopsignal=TERM
stopasgroup=true
killasgroup=true

# todo: maybe a wrapper with traps are better here
# currently postfix is started in entrypoint.sh
[program:postfix]
priority=5
command=/usr/sbin/postfix -c /etc/postfix start
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
startsecs=10
autorestart=true
autostart=false
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:postgresql]
priority=5
command=/usr/bin/postgres -D /opt/database
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
user=postgres
startretries=10
startsecs=10
autorestart=true
autostart=false
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:mosquitto]
priority=5
command=mosquitto -c /etc/mosquitto.conf
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
user=root
startretries=10
startsecs=2
autorestart=true
autostart=false
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:ospd-openvas]
command=sudo -Hiu gvm /opt/setup/scripts/ospd-openvas-wrapper.sh /usr/local/bin/ospd-openvas -f --unix-socket /var/run/ospd/ospd.sock --pid-file /var/run/ospd/ospd.pid --socket-mode 0o777 --log-level INFO
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
priority=10
startretries=10
startsecs=10
autorestart=true
autostart=false
depends_on=redis,postgresql,mosquitto
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:gvmd]
command=/opt/setup/scripts/start_gvmd.sh -f --osp-vt-update=/var/run/ospd/ospd.sock --listen=0.0.0.0 --port=9390 --gnutls-priorities=SECURE128:-AES-128-CBC:-CAMELLIA-128-CBC:-VERS-SSL3.0:-VERS-TLS1.0:-VERS-TLS1.1
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
priority=20
user=gvm
startretries=30
startsecs=10
autorestart=true
autostart=false
depends_on=redis,postgresql,ospd-openvas
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:gsad-https-owncert]
command=/opt/setup/scripts/start_gsad.sh -f --verbose --gnutls-priorities=SECURE128:-AES-128-CBC:-CAMELLIA-128-CBC:-VERS-SSL3.0:-VERS-TLS1.0:-VERS-TLS1.1 --ssl-certificate=%(ENV_CERTIFICATE)s --ssl-private-key=%(ENV_CERTIFICATE_KEY)s --timeout=%(ENV_TIMEOUT)s --no-redirect --mlisten=127.0.0.1 --mport=9390 --port=9392
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
priority=30
user=gvm
startretries=5
startsecs=10
autorestart=true
autostart=false
depends_on=redis,postgresql,ospd-openvas,gvmd
stopsignal=KILL
stopasgroup=true
killasgroup=true

[program:gsad-https]
command=/opt/setup/scripts/start_gsad.sh -f --verbose --gnutls-priorities=SECURE128:-AES-128-CBC:-CAMELLIA-128-CBC:-VERS-SSL3.0:-VERS-TLS1.0:-VERS-TLS1.1 --timeout=%(ENV_TIMEOUT)s --no-redirect --mlisten=127.0.0.1 --mport=9390 --port=9392
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
priority=30
user=gvm
startretries=5
startsecs=10
autorestart=true
autostart=false
depends_on=redis,postgresql,ospd-openvas,gvmd
stopsignal=KILL
stopasgroup=true
killasgroup=true

[program:gsad]
command=/opt/setup/scripts/start_gsad.sh -f --verbose --http-only --timeout=%(ENV_TIMEOUT)s --no-redirect --mlisten=127.0.0.1 --mport=9390 --port=9392
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_err.log
priority=30
user=gvm
startretries=5
startsecs=10
autorestart=true
autostart=false
depends_on=redis,postgresql,ospd-openvas,gvmd
stopsignal=KILL
stopasgroup=true
killasgroup=true

#28800 = 8h
[program:GVMUpdate]
command=bash -c "sleep 28800 && exec /opt/setup/scripts/sync-all.sh"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=900
user=root
startretries=10
startsecs=10
autorestart=true
autostart=false
depends_on=redis,postgresql,ospd-openvas,gvmd
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:crond]
command=/opt/setup/scripts/start_crond.sh
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=900
startretries=10
startsecs=10
autorestart=true
autostart=true
stopsignal=TERM
stopasgroup=true
killasgroup=true

